# postgresql12.service
# Adapted from https://www.postgresql.org/docs/devel/static/server-start.html
[Unit]
Description=PostgreSQL database server
Documentation=man:postgres(1)
After=postgresql12-install.service
Requires=postgresql12-install.service

[Service]
Type=notify
User=postgres
ExecStart=/usr/libexec/postgresql12/postgres -D /var/lib/pgsql12/data
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGINT
TimeoutSec=0

[Install]
WantedBy=multi-user.target
